# package/linux/pacta.service
[Unit]
Description=PACTA Contract Management System
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pacta
ExecStart=/opt/pacta/node/bin/node /opt/pacta/server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
